# ai_service.ini
[program:ai_service]
numprocs=1
process_name=%(program_name)s_%(process_num)02d;
command=/Users/humin/Documents/virtual/ai_service/bin/python3 /Users/humin/Documents/virtual/ai_service/bin/gunicorn -c /Users/humin/Documents/project/ai_service/gunicorn_config.py app:app
directory=/Users/humin/Documents/project/ai_service/
user=humin
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true